# memory-reclaim/mem-reclaim.timer
[Unit]
Description=Aggressive memory reclaim script for WSL2

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target